# scripts/app-registration.service
[Unit]
Description=APP/小程序登记管理系统
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/app-registration-system
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=9999

# 日志
StandardOutput=append:/var/log/app-registration-system/out.log
StandardError=append:/var/log/app-registration-system/error.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
